# Dockerfile/move2docker/etc/service/code-server/run
#!/bin/bash
su $USERNAME -c "/etc/code-server-hub/.cshub/code-server --host 127.0.0.1 --port 16313 --auth none"
